# shutdown-hook.conf
description "run at shutdown"

start on starting rc
task
exec /bin/bash /home/ubuntu/reverse_proxy/app/utils/sync_rotated_file_to_s3.sh log
